# resources/cron/cloud_ops_esmachine.cron
3,13,23,33,43,53 * * * * /home/klahnakoski/esReplicate/resources/scripts/replicate_orange_factor.sh
